# etc/supervisord.d/nfs-mounts.conf
[program:nfs-mounts]
command=/opt/bin/nfs.mounts
priority=200
exitcodes=0
stdout_events_enabled=true
stderr_events_enabled=true
