# etc/supervisord.d/objectivefs.conf
[program:objectivefs]
command=mount.objectivefs mount %(ENV_BUCKET)s /objectivefs
exitcodes=0
stdout_events_enabled=true
stderr_events_enabled=true
